# config/systemd/user/low-battery-notifier.service
[Unit]
Description=Check battery level and send notification if low

[Service]
Type=oneshot
ExecStart=%h/.scripts/battery/low-battery-alert.sh
