# config/systemd/user/mnt-qnap.automount
[Unit]
Description=Automount Qnap

[Automount]
Where=/mnt/qnap

[Install]
WantedBy=multi-user.target
